# onceaday/crontab
@daily python3 /app/onceaday/onceaday.py /app/onceaday
